# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1280x800x24 -ac
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null
priority=10

[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -shared -rfbport 5900 -passwd %(ENV_VNC_PASSWORD)s
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null
depends_on=xvfb
priority=20

[program:novnc]
command=/usr/share/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 6080
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null
depends_on=x11vnc
priority=30

[program:app]
command=uvicorn server:app --host 0.0.0.0 --port 8900
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=DISPLAY=":99",PYTHONUNBUFFERED="1"
priority=40
